# net-misc/notif/files/notif.service
[Unit]
Description=notif server

[Service]
ExecStart=/usr/bin/notif route
User=nobody
Restart=always

[Install]
WantedBy=multi-user.target
